# templates/well-firstboot.service
[Unit]
Description=Wells per-well first-boot identity injection
Documentation=https://github.com/peteknowsai/splites/blob/main/docs/MVP-PLAN.md#b09d4
# Idempotent: subsequent boots skip the whole unit in microseconds.
# /etc/.well-ready is touched at the end of well-firstboot, after
# identity is fully applied.
ConditionPathExists=!/etc/.well-ready
# Run AFTER networkd is up. This is the inverse of the earlier B-shape
# (Before=systemd-networkd.service) which was rejected: ordering before
# networkd hung the entire boot when the cidata seed was missing on a
# steady-state restart, because well-firstboot's failure path blocked
# downstream networking. Late hostname is fine — welld locates new
# wells via delta-snapshot lease lookup, not hostname matching.
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/well-firstboot
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
